# packaging/memos-rs.service
[Unit]
Description=memos-rs
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=memos-rs
Group=memos-rs
WorkingDirectory=/var/lib/memos-rs
ExecStart=/usr/local/bin/memos-rs serve --config /etc/memos-rs/config.toml
Restart=on-failure
RestartSec=5
StateDirectory=memos-rs
StateDirectoryMode=0750
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native
ReadWritePaths=/var/lib/memos-rs

[Install]
WantedBy=multi-user.target
